# sv/gitaly/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

exec 2>&1
exec sudo -u git -H /home/git/gitaly/gitaly /home/git/gitaly/config.toml
